# config/supervisor/worker.conf
; ============================================================
; Supervisor Configuration — Queue Worker
; ============================================================
; Manages Laravel queue worker processes
; Config via .env: WORKER_PROJECT, WORKER_QUEUE_CONNECTION
;
; Commands:
;   docker compose exec php-worker supervisorctl status
;   docker compose exec php-worker supervisorctl restart all
;   docker compose exec php-worker supervisorctl stop queue-worker:*
; ============================================================

[supervisord]
nodaemon=true
user=root
logfile=/var/log/worker/supervisord.log
pidfile=/var/run/supervisord.pid
logfile_maxbytes=10MB

[program:queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/%(ENV_WORKER_PROJECT)s/artisan queue:work %(ENV_WORKER_QUEUE_CONNECTION)s --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/worker/queue-%(process_num)02d.log
stdout_logfile_maxbytes=10MB
stopwaitsecs=3600
startsecs=0
